# supervisor/conf.d/apache.conf
[program:apache]
command=apache2ctl -D FOREGROUND
redirect_stderr=true
