# supervisor/conf.d/gulden.conf
[program:gulden]
command=/opt/gulden/gulden/GuldenD -datadir=/opt/gulden/datadir
redirect_stderr=true
